# epoint/files/etc/xinetd.d/epoint
service epoint
{
    disable     = no
    protocol    = tcp
    port        = 8080
    socket_type = stream
    wait        = no
    server      = /usr/lib/epoint/http-redirect.awk
    user        = root
    group       = root
}
